# envoy.service
[Service]
ExecStart=/usr/bin/envoyd

[Install]
Also=envoy.socket
