# envoy.socket
[Socket]
ListenStream=/run/keychain

[Install]
WantedBy=sockets.target
