# envoy@.service
[Service]
ExecStart=/usr/bin/envoyd --agent=%i

[Install]
Also=envoy@%i.socket
