# envoy.socket
[Socket]
ListenStream=@/vodik/envoy

[Install]
WantedBy=sockets.target
