# systemd/envoy.service
[Unit]
Description=Envoy agent monitor

[Service]
ExecStart=/usr/bin/envoyd
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
Also=envoy.socket
